# root/etc/s6-overlay/s6-rc.d/svc-mod-universal-openvpn/run
#!/usr/bin/with-contenv bash
if [[ ${OVPN:-0} -eq 1 ]]; then
    exec s6-setuidgid abc openvpn --config ${OVPN_CONFIG:-/etc/openvpn/config.ovpn} --auth-user-pass ${OVPN_CREDENTIALS:-/etc/openvpn/credentials.txt}
else
    exec s6-setuidgid abc bash -c -- 'sleep 15; while true; do sleep infinity; done;'
fi