# root/etc/s6-overlay/s6-rc.d/init-mod-universal-openvpn/run
#!/usr/bin/with-contenv bash

echo "openvpn" >> /mod-repo-packages-to-install.list